# tests/scripts/kdump-test.service
[Unit]
Description=Kdump Test Service

[Service]
ExecStart=/test-entry.sh
Type=idle

[Install]
WantedBy=multi-user.target
